# scripts/crontab
*/5 * * * * /home/wise/data/SITE-CHECK/site-check/scripts/runonce.py


